# conf/systemd.service
[Unit]
Description=Etherpad-lite, the collaborative editor.
After=syslog.target network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
Environment="PATH=__ENV_PATH__"
Environment=NODE_ENV=production
ExecStart=__FINALPATH__/src/bin/run.sh /var/log/__APP__/etherpad.log
Restart=always

[Install]
WantedBy=multi-user.target
